# ansible_project_sourcecode/update_sop/files/usr/lib/systemd/system/ansible_send_mail.service
[Unit]
Description=M

[Service]
Type=simple
ExecStart=/home/ansible/SOP_updater/run_sop.sh report-mail

[Install]
WantedBy=multi-user.target